# dotfiles/systemd/user/post-sleep.target
[Unit]
Description=User level post-sleep target
StopWhenUnneeded=yes
