# dotfiles/systemd/user/pre-sleep.target
[Unit]
Description=User level pre-sleep target
StopWhenUnneeded=yes
