# dotfiles/picom/stop-picom.service
[Service]
ExecStart=/usr/bin/systemctl --user stop picom
Type=exec

[Install]
WantedBy=pre-sleep.target
